# src/supervisord.conf
[program:monitor]
command = dotnet exec /app/dotnet-monitor.dll collect --no-auth
stdout_logfile = /dev/stdout
stderr_logfile = /dev/stderr

[inet_http_server]
port=127.0.0.1:20000
